# .config/systemd/user/desktop-online.target
[Unit]
Description=Desktop is up and running
After=desktop.target
